# Capture/filesystem/etc/supervisord.conf
[supervisord]
autostart=true
autorestart=true
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor
pidfile=/var/run/supervisord.pid

[program:nfcapd]
command = nfcapd_start.sh
autostart = true
autorestart = true

[program:yaf]
command = yaf_start.sh
autostart = true
autorestart = true

[program:cron] 
command = cron -f -L 15
autostart = true
autorestart = true
